# dyndns.timer
[Unit]
Description=Run DynDNS client every 10 minutes

[Timer]
OnBootSec=1min
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
